# docker/supervisor.conf
[program:janitor]
command=/usr/local/bin/gunicorn -b localhost:8000 -w 4 janitor:app --preload
directory=/opt/janitor
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
